# root/etc/services.d/flexget/run
#!/usr/bin/with-contenv bash

s6-setuidgid abc /usr/bin/flexget -c /app/all_flows.yml -l /config/logs/flexget.log daemon start
